# .config/extend-rc/autostart/clipit.desktop
[Desktop Entry]
Exec=pgrep -x clipit >/dev/null && clipit
Name=clipit
Type=Application
