# .config/extend-rc/autostart/optimus-manager-qt.desktop
[Desktop Entry]
Exec=QT_QPA_PLATFORMTHEME=qt5ct /usr/bin/optimus-manager-qt
Name=optimus-manager-qt
Type=Application
