# .config/extend-rc/autostart/pa-applet.desktop
[Desktop Entry]
Exec=/usr/bin/pa-applet
Name=pa-applet
Type=Application
