# .config/extend-rc/autostart/picom.desktop
[Desktop Entry]
Exec=/usr/bin/picom -b
Name=picom
Comment=a compositor for X11
Type=Application
