# .config/extend-rc/autostart/polybar.desktop
[Desktop Entry]
Exec=$HOME/.config/polybar/launch.sh
Name=launch.sh
Type=Application
